# somel-apache2/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Démarrage d'Apache2..."
exec /usr/sbin/httpd -D FOREGROUND
